# offlineimap.service
[Unit]
Description=Offlineimap Service

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap -o
User=adam

[Install]
WantedBy=mail.target
