# systemd/VBoxClient.service
[Unit]
Description=VirtualBox helpers

[Service]
ExecStart=/usr/bin/VBoxClient-all
Type=forking

[Install]
WantedBy=graphical-session.target